# Configs/.config/systemd/user/zsh-zcompdump-clean.service
[Unit]
Description=Clean stale zcompdump temp files

[Service]
Type=oneshot
ExecStart=%h/.local/bin/zsh-clean-zcompdump
