# Configs/.config/sv/power-profile-auto/run
#!/bin/sh
# runit service: automatic power-profile switching.
# Idle (rather than crash-loop) when power-profiles-daemon is unavailable.
exec 2>&1
if ! command -v powerprofilesctl >/dev/null 2>&1; then
	echo "power-profile-auto: powerprofilesctl not found; idling"
	exec sleep infinity
fi
exec "$HOME/.local/lib/hypr/system/powerprofiles.auto.sh"
