# systemd/.config/systemd/user/hypr-config.service
[Unit]
Description=Hyprland Configuration Parser Service (TOML to ENV)
Documentation=https://github.com/mimabial/hyprdots
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/state/hypr/pip_env/bin/python %h/.local/lib/hypr/util/parse.config.py --daemon --export
Restart=on-failure
RestartSec=5s
Environment="DISPLAY=:0"
Environment="PYTHONPATH=%h/.local/lib/hypr"

# Make sure the required directories exist
ExecStartPre=/bin/mkdir -p %h/.config/hypr
ExecStartPre=/bin/mkdir -p %h/.local/state/hypr

[Install]
WantedBy=graphical-session.target
