# systemd/.config/systemd/user/hyprland-keep-awake.service
[Unit]
Description=Hyprland keep-awake inhibitor
Documentation=https://wiki.hyprland.org/Hypr-Ecosystem/hypridle
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStart=/usr/bin/systemd-inhibit --what=idle --who=hyprland-keep-awake --why="Keep awake enabled" --mode=block /usr/bin/sleep infinity
Restart=on-failure

[Install]
WantedBy=graphical-session.target
