# systemd/.config/systemd/user/inhibit-idle-on-audio.service
[Unit]
Description=Inhibit idle while audio plays
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/inhibit-idle-on-audio.sh
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
