# systemd/.config/systemd/user/zsh-zcompdump-clean.timer
[Unit]
Description=Daily cleanup of zcompdump temp files

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
