# Configs/.config/sv/auto-theme/run
#!/bin/sh
# runit service: automatic light/dark theme daemon.
# Down by default — color-mode brings it up/down (sv up/down auto-theme) when
# the colour mode is set to "auto". See the `down` file in this directory.
exec 2>&1
exec "$HOME/.local/state/hypr/pip_env/bin/python" \
	"$HOME/.local/lib/hypr/theme/auto_theme.py"
