# Configs/.config/sv/hypr-config/run
#!/bin/sh
# runit service: Hyprland config parse/export daemon.
exec 2>&1
exec "$HOME/.local/state/hypr/pip_env/bin/python" \
	"$HOME/.local/lib/hypr/util/parse.config.py" --daemon --export
